# conf/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:ams]
command=/etc/init.d/ams start
